# examples/shutdownd.bsdrc
#!/bin/sh
. /etc/rc.subr

name="shutdownd"
rcvar=shutdownd_enable
pidfile="/var/run/${name}.pid"
command="/usr/local/bin/shutdownd"

load_rc_config ${name}
run_rc_command "$1"
